# deploy/claude-token-proxy.service
[Unit]
Description=Claude Token Proxy Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/nexus/infra/claude-token-proxy
EnvironmentFile=%h/nexus/infra/claude-token-proxy/.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
